# packaging/out/deepcool-desktop_0.1.0_amd64/etc/systemd/system/deepcool-digital-dart.service
[Unit]
Description=DeepCool Digital Dart Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/deepcool-digital-dart
Restart=on-failure

[Install]
WantedBy=multi-user.target
