# packaging/udev/99-deepcool-digital.rules
# udev rules for DeepCool Digital HID devices
KERNEL=="hidraw*", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="3633", MODE:="0666", TAG+="uaccess"
KERNEL=="hidraw*", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="34d3", ATTRS{idProduct}=="1100", MODE:="0666", TAG+="uaccess"

# Let the desktop app read CPU package power from Linux powercap/RAPL.
SUBSYSTEM=="powercap", KERNEL=="*-rapl:*", RUN+="/bin/chmod a+r /sys/class/powercap/%k/energy_uj", RUN+="/bin/chmod a+r /sys/class/powercap/%k/max_energy_range_uj"
